# conf/hello_app/supervisor.conf
[program:hello]
command=/home/box/web/apps/hello/gunicorn_start.sh
stdout_logfile=/home/box/web/apps/hello/log/supervisor.log
redirect_stderr=true