# conf/ask_app/supervisor_ask.conf
[program:ask]
command=/home/box/web/ask/gunicorn_start.sh
stdout_logfile=/home/box/web/ask/log/supervisor.log
redirect_stderr=true